# self.service
[Unit]
Description=example systemd service unit file.

[Service]
ExecStart=/bin/bash /usr/sbin/self.sh

[Install]
WantedBy=multi-user.target